# deploy/systemd/iou-ai-event-projector.timer
[Unit]
Description=Check for redacted io_uring notification events

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=iou-ai-event-projector.service

[Install]
WantedBy=timers.target

