# deploy/systemd/iou-ai-export.timer
[Unit]
Description=Hourly read-only AFL++ fleet-stat export

[Timer]
OnCalendar=hourly
Persistent=true
Unit=iou-ai-export.service

[Install]
WantedBy=timers.target

